# files/cron.d/cron-apt
#
# Regular cron jobs for the cron-apt package
#
# Every night at 4 o'clock.
0 4	* * <%= $dotw %>	root	test -x /usr/sbin/cron-apt && /usr/sbin/cron-apt
